# observability/grafana/cron/report-filespace.crontab
* * * * * export BYTES=`df -P /dev/vda1 | awk 'NR==2 {print $4}'`; curl -k -X POST "https://telemetry.hackclub.com/mtx?k=production.grafana.free_disk_space_bytes&t=g&v=$BYTES" --user hackclub:replace_me_pw --verbose

